# init.ini
[program:mariadb]
command=/usr/bin/mysqld --user=mysql
user=mysql
startsecs=5
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_err.log
priority=110

[program:redis]
command=/usr/bin/redis-server /etc/redis.conf --sentinel
user=redis
startsecs=5
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_err.log
priority=111

[program:memcached]
command=/usr/bin/memcached -m 64 -p 11211 -u memcached -l 127.0.0.1
user=memcached
startsecs=5
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_err.log
priority=112

[program:apache]
command=/usr/sbin/apachectl -DFOREGROUND
user=root
startsecs=5
autostart=true
autorestart=true
killasgroup=true
stopasgroup=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_err.log
priority=113

[program:php-fpm7]
command=/usr/sbin/php-fpm7 -F
user=root
startsecs=5
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_err.log
priority=114
